# roles/absvfx.workstation/templates/ansible-network.service.j2
[Unit]
Description=Ansible Temporary Network

[Service]
ExecStart=/bin/bash /usr/local/bin/ansible-network.sh

[Install]
WantedBy=multi-user.target
